# target-linux/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:logdir]
command=/bin/bash -c "mkdir -p /var/log/sensel-training && chown sensel-svc:sensel-svc /var/log/sensel-training && sleep infinity"
autostart=true
autorestart=true
priority=10
user=root

[program:sandcat]
command=/opt/sensel/scripts/bootstrap-sandcat.sh
autostart=true
autorestart=true
startsecs=5
startretries=10
priority=20
user=root
stdout_logfile=/var/log/sensel-training/sandcat.log
stderr_logfile=/var/log/sensel-training/sandcat.err.log

[program:wazuh-agent]
command=/bin/bash -c "if [ \"${ENABLE_WAZUH:-false}\" = \"true\" ]; then exec /var/ossec/bin/wazuh-control start; else sleep infinity; fi"
autostart=true
autorestart=true
priority=30
user=root
stdout_logfile=/var/log/sensel-training/wazuh.log
stderr_logfile=/var/log/sensel-training/wazuh.err.log
